# udev-rules/10-input-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/victor/bin/monitor-switcher"
